# scripts/com.infiniterecall.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
        "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.infiniterecall.api</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/infinite-recall-api</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info,infinite_recall_api=info</string>
        <key>INFINITE_RECALL_BIND</key>
        <string>127.0.0.1:7331</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/infinite-recall-api.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/infinite-recall-api.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
